# systemd/partle-mcp.service
[Unit]
Description=Partle MCP Servers for ChatGPT Integration
After=network.target partle-backend.service
Wants=partle-backend.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/partle/backend
Environment=PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/root/partle/backend
Environment=PARTLE_API_URL=http://localhost:8000
ExecStart=/root/.local/bin/uv run python scripts/run_mcp_products.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=partle-mcp

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/partle /var/log/partle /tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target