# systemd/partle-backend.service
[Unit]
Description=Partle Backend FastAPI Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/srv/partle/backend
Environment=PATH=/home/deploy/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/srv/partle/backend
ExecStartPre=/home/deploy/.local/bin/poetry check
ExecStart=/home/deploy/.local/bin/poetry run uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=partle-backend

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/partle /var/log/partle /tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target